# runit/keymux-hyprland/run
#!/bin/sh
exec /usr/bin/keymux hyprland-daemon 2>&1
