# keyboard-middleware.service
[Unit]
Description=Keyboard Middleware Daemon
Documentation=https://github.com/yourusername/keyboard-middleware
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/keyboard-middleware daemon
Restart=on-failure
RestartSec=5
Environment="RUST_LOG=debug"

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=default.target
