# runit/keymux-sway/run
#!/bin/sh
exec /usr/bin/keymux sway-daemon 2>&1
