# runit/keymux/run
#!/bin/sh
exec /usr/bin/keymux daemon 2>&1
