# systemd/keymux-hyprland.service
[Unit]
Description=Keyboard Middleware Hyprland Window Watcher
Documentation=https://github.com/fibsussy/keymux
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/keymux hyprland-daemon
Restart=on-failure
RestartSec=5
Environment="RUST_LOG=info"

# Security hardening
# NOTE: options that create mount namespaces (ProtectSystem, ProtectHome,
# ProtectKernelTunables, ProtectKernelModules, ProtectControlGroups, etc.)
# block /proc/{pid}/environ reads needed for IS_GAME=1 game mode detection.
RestrictRealtime=true

[Install]
WantedBy=graphical-session.target
