# runit/keymux-niri/run
#!/bin/sh
exec /usr/bin/keymux niri-daemon 2>&1
