# openrc/keymux
#!/sbin/openrc-run

name="keymux"
description="Keyboard Middleware Async Root Daemon"
command="/usr/bin/keymux"
command_args="daemon"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/keymux.log"
error_log="/var/log/keymux.err"

depend() {
    need localmount
    after firewall
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /run
    checkpath --directory --owner root:root --mode 0755 /var/log
}

start() {
    ebegin "Starting ${name}"
    start-stop-daemon --start --exec ${command} -- ${command_args} \
        --pidfile ${pidfile} \
        --stdout ${output_log} \
        --stderr ${error_log}
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop --exec ${command} --pidfile ${pidfile}
    eend $?
}

reload() {
    ebegin "Reloading ${name} configuration"
    /usr/bin/keymux reload
    eend $?
}
