# runit/keymux-bspwm/run
#!/bin/sh
exec /usr/bin/keymux bspwm-daemon 2>&1
