# runit/keymux-i3/run
#!/bin/sh
exec /usr/bin/keymux i3-daemon 2>&1
